# homesolar/rootfs/etc/s6-overlay/s6-rc.d/homesolar/run
#!/usr/bin/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Home Assistant Add-on: HomeSolar
# Runs the HomeSolar application
# ==============================================================================

declare latitude
declare longitude
declare timezone
declare auto_detect_language
declare ha_language

# Read configuration
latitude=$(bashio::config 'latitude')
longitude=$(bashio::config 'longitude')
timezone=$(bashio::config 'timezone')
auto_detect_language=$(bashio::config 'auto_detect_language')

# Get ingress settings
export INGRESS_ENTRY=$(bashio::addon.ingress_entry)

# Get Home Assistant language
ha_language=$(bashio::info.language)

# Export variables for the app
export LATITUDE="${latitude}"
export LONGITUDE="${longitude}"
export TIMEZONE="${timezone}"
export AUTO_DETECT_LANGUAGE="${auto_detect_language}"
export HA_LANGUAGE="${ha_language}"

bashio::log.info "Starting HomeSolar - Solar Ephemeris..."
bashio::log.info "Position: ${latitude}, ${longitude}"
bashio::log.info "Timezone: ${timezone}"
bashio::log.info "Language: ${ha_language}"

cd /app || bashio::exit.nok "Cannot change to /app directory"
exec gunicorn --bind 0.0.0.0:8099 --workers 1 --threads 2 app:app
